# etc/profile.d/elm.sh
export ELM_PROFILE="slp"
export ELM_MODULES="ctxpopup_copypasteshareUI>entry/api:decarta>map/api"
