# etc/profile.d/elm.sh
export ELM_FONT_PATH=/usr/share/SLP/fonts
export ELM_PAGE_SCROLL_FRICTION=0.5
export ELM_BRING_IN_SCROLL_FRICTION=0.5
export ELM_FPS=100
export ELM_INPUT_PANEL=1
