# etc/profile.d/elm.sh
export ELM_FONT_PATH=/usr/share/SLP/fonts
export ELM_PAGE_SCROLL_FRICTION=0.5
export ELM_THUMBSCROLL_FRICTION=3.0
export ELM_BRING_IN_SCROLL_FRICTION=0.5
export ELM_FPS=6000
export ELM_INPUT_PANEL=1
export EDJE_MULTISENSE_UI_MODULE=slp_sound_haptic
export ELM_MODULES="ctxpopup_copypasteUI>entry/api:decarta>map/api"
export ELM_PASSWORD_SHOW_LAST_CHARACTER=1
export ELM_LONGPRESS_TIMEOUT=0.5
