# etc/profile.d/elm.sh
export ELM_PROFILE="slp"
export ELM_MODULES="ctxpopup_copypasteUI>entry/api:decarta>map/api"
