# etc/profile.d/elm.sh
export ELM_MODULES="ctxpopup_copypasteUI>entry/api:datetime_input_ctxpopup>datetime/api"
# efl-1.8
export ELM_MODULES="$ELM_MODULES:prefs>prefs_iface"

export ELM_THEME="tizen-HD-light"
